# homedir/.crontab
# ruby gem nightly updates
0 2 * * *  gem update -V -y
# brew nightly updates
0 3 * * *  brew update
0 4 * * *  brew upgrade
